# install/cron/www.cron
GREENERY_WEB=/var/www/greenery

* * * * * $GREENERY_WEB/bin/scripts/scheduler.py
* * * * * $GREENERY_WEB/bin/scripts/actions.py
* * * * * $GREENERY_WEB/bin/scripts/poll.py
0 * * * * $GREENERY/bin/scripts/aggregate.py
