# install/cron/www.cron
GREENERY_WEB=/var/www/greenery

# check for scheduled actions
* * * * * $GREENERY_WEB/greenery/bin/scripts/scheduler.py

# check for required threshold actions
* * * * * $GREENERY_WEB/greenery/bin/scripts/actions.py

# poll sensors for data
* * * * * $GREENERY_WEB/greenery/bin/scripts/poll.py

# hourly aggregate previous hours measurement data
0 * * * * $GREENERY_WEB/greenery/bin/scripts/aggregate.py
