# install/cron/www.cron
GREENERY_WEB=/var/www/greenery

# poll sensors for data
* * * * * $GREENERY_WEB/greenery/bin/poll.py

# check for scheduled actions
* * * * * $GREENERY_WEB/greenery/bin/scheduler.py

# check for required threshold actions
* * * * * $GREENERY_WEB/greenery/bin/actions.py

# aggregate previous hours measurement data
0 * * * * $GREENERY_WEB/greenery/bin/aggregate.py

# purge expired measurement data
5 * * * * $GREENERY_WEB/greenery/bin/purge.py

