# conf/www.cron
GREENERY_WEB=/var/www/greenery

# scheduler for greenery.guru app
* * * * * $GREENERY_WEB/scripts/scheduler.py
